# .config/systemd/user/low-battery-notifier.service
[Unit]
Description=Notify low battery

[Service]
Type=oneshot
ExecStart=/home/bruhtus/.local/bin/low-battery-notifier
