# .config/systemd/user/updates-notifier.service
[Unit]
Description=Check and notify pacman and AUR updates

[Service]
Type=simple
ExecStart=/home/bruhtus/.local/bin/updates-notifier
